# roles/security/templates/lynis.timer.j2
[Unit]
Description=Server Helper Lynis Audit Timer
Requires=server-helper-lynis.service

[Timer]
OnCalendar={{ security.lynis.schedule }}
Persistent=true

[Install]
WantedBy=timers.target
