# roles/restic/templates/restic-backup.service.j2
# Restic Backup Service
# Generated by Server Helper Ansible

[Unit]
Description=Restic Backup Service
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/scripts/restic-backup.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=restic-backup

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
