# roles/restic/templates/restic-backup.timer.j2
# Restic Backup Timer
# Generated by Server Helper Ansible
# Schedule: {{ restic.schedule }}

[Unit]
Description=Restic Backup Timer
Requires=restic-backup.service

[Timer]
# Cron-style schedule
OnCalendar={{ restic.schedule }}

# Run on boot if missed
Persistent=true

# Randomize start time by up to 10 minutes
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
