# roles/lynis/templates/lynis-scan.service.j2
[Unit]
Description=Lynis Security Scan
Documentation=https://cisofy.com/lynis/
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/lynis-scan.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lynis-scan

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes

# Timeout (scan can take a while)
TimeoutStartSec=10min

[Install]
WantedBy=multi-user.target
