# roles/netdata/templates/netdata-health-check.service.j2
[Unit]
Description=Netdata Health Check
Documentation=https://github.com/thelasttenno/Server-Helper
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/check_netdata_health.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=netdata-health-check

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes
