# roles/backups/templates/backup.timer.j2
[Unit]
Description=Server Helper Backup Timer
Requires=server-helper-backup.service

[Timer]
OnCalendar={{ backups.schedule }}
Persistent=true

[Install]
WantedBy=timers.target
